# systemd/user-units/dbus-system-i3blocks-notifier.service
[Unit]
Description=Notify i3blocks when certain D-Bus events occur

[Service]
ExecStart=%h/.dotfiles/systemd/scripts/dbus-system-i3blocks-notifier.sh